# Scripts/Udev/Test-Scripts/Udev/Rules/old/30-usb-drive_by-defender.rules
ACTION=="add", ATTRS{idVendor}=="****", ATTRS{idProduct}=="****", RUN+="/usr/local/bin/my_script.sh"
